# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/nginx.stdout.log
stderr_logfile=/var/log/nginx.stderr.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/php-fpm.stdout.log
stderr_logfile=/var/log/php-fpm.stderr.log
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
